# service.service
[Service]
ExecStart=/home/vp347/project/Project/Front/dbLister.php
[Install]
WantedBy=multi-user.target
